# systemd/audiobook-shutdown-saver.service
[Unit]
Description=Audiobook Staging Saver (save tmpfs before shutdown)
Documentation=https://github.com/TheBoscoClub/Audiobook-Manager
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
# audiobook-save-staging enforces audiobooks user via require_audiobooks_user()
# so newly-moved files (staging → library) inherit the correct audiobooks:audiobooks
# ownership. Running as root would create root-owned files in the library and
# break the project's permission model. The script's `systemctl stop` calls for
# audiobook-converter/mover are silently ignored when invoked without privilege —
# acceptable because the shutdown sequence stops them shortly afterwards anyway.
# (Audiobook-Manager-c8l: prior `User=root` caused every shutdown attempt to FAIL
# with "audiobook-save-staging must run as the audiobooks user" — staging never
# persisted on reboot.)
User=audiobooks
Group=audiobooks

# Load configuration
EnvironmentFile=-/etc/audiobooks/audiobooks.conf

# Handle stdout at shutdown (journal may be unavailable)
StandardOutput=journal
StandardError=journal

# Save staging files before shutdown
ExecStart=/usr/local/bin/audiobook-save-staging

# Run only at shutdown
RemainAfterExit=yes

# Security hardening (PrivateTmp intentionally omitted — needs real /tmp to save staging files)
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=halt.target reboot.target shutdown.target
