# systemd/audiobook-stream-translate.service
[Unit]
Description=Audiobook streaming translation worker
Documentation=https://github.com/TheBoscoClub/Audiobook-Manager
After=network-online.target audiobook-api.service
Wants=network-online.target
PartOf=audiobook.target
# Ensure application tree and data directory are mounted (RAID arrays can be
# slow at boot). Matches audiobook-api.service convention.
RequiresMountsFor=/opt/audiobooks /srv/audiobooks
# Generous interval caps restart storms if the worker is persistently broken
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=audiobooks
Group=audiobooks

# Do NOT set EnvironmentFile=/etc/audiobooks/audiobooks.conf here.
#
# audiobooks.conf uses shell-assignment syntax with ${var} references — e.g.
# AUDIOBOOKS_LIBRARY="${AUDIOBOOKS_DATA}/Library". systemd's EnvironmentFile
# parser reads key=value only: it does NOT perform shell expansion, so
# AUDIOBOOKS_LIBRARY would be injected into the child environment as the
# literal string "${AUDIOBOOKS_DATA}/Library". Our wrapper then sources
# audiobook-config.sh which uses `: "${VAR:=default}"` (assign-only-if-unset),
# so the pre-set literal wins and the worker fails with "library directory
# not found: ${AUDIOBOOKS_DATA}/Library".
#
# The wrapper (stream-translate-daemon.sh) sources lib/audiobook-config.sh
# directly, which does do shell expansion correctly — that is the canonical
# path.

WorkingDirectory=/opt/audiobooks

# Wrapper resolves AUDIOBOOKS_DATABASE / AUDIOBOOKS_LIBRARY / AUDIOBOOKS_VENV
# from the canonical bash config loader (not from systemd env expansion) and
# execs scripts/stream-translate-worker.py.
ExecStart=/opt/audiobooks/scripts/stream-translate-daemon.sh

# Restart on failure; StartLimit* caps (in [Unit]) bound restart storms
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audiobook-stream-translate

# Security hardening — matches house style used by audiobook-converter.service.
# ProtectSystem=full write-protects /usr,
# /boot, /etc but leaves /var, /srv, and /tmp writable by default, which is
# what the worker needs for the DB, the audiobook library, and per-segment
# opus audio under AUDIOBOOKS_STREAMING_AUDIO_DIR. We deliberately do NOT use
# ProtectSystem=strict because it requires an explicit ReadWritePaths list,
# and systemd does not expand EnvironmentFile vars in ReadWritePaths — any
# user override of AUDIOBOOKS_VAR_DIR or AUDIOBOOKS_LIBRARY would silently
# break the sandbox.
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictNamespaces=yes
SystemCallArchitectures=native

[Install]
WantedBy=audiobook.target
