# systemd/audiobook-downloader.timer
[Unit]
Description=Audiobook Download Timer (check for new audiobooks)
Documentation=https://github.com/TheBoscoClub/Audiobook-Manager
PartOf=audiobook.target

[Timer]
# Run 2 minutes after boot
OnBootSec=2min

# Then run every 5 minutes
OnUnitActiveSec=5min

# Persist timer across reboots
Persistent=true

[Install]
WantedBy=audiobook.target timers.target
