# systemd/audiobook-translate-check.service
[Unit]
Description=Check for pending audiobook translations
Documentation=https://github.com/TheBoscoClub/Audiobook-Manager
After=network-online.target audiobook-api.service

[Service]
Type=oneshot
ExecStart=/opt/audiobooks/scripts/translation-check.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audiobook-translate-check

# Security hardening
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
