# systemd/audiobook.target
[Unit]
Description=Audiobook Library Services
Documentation=https://github.com/TheBoscoClub/Audiobook-Manager
After=network-online.target
Wants=network-online.target

# Core services - these are the main services that must start
Wants=audiobook-api.service
Wants=audiobook-proxy.service
Wants=audiobook-redirect.service

# Background workers - conversion and file moving
Wants=audiobook-converter.service
Wants=audiobook-mover.service
Wants=audiobook-stream-translate.service

# Scheduled tasks - timers trigger their corresponding .service units
Wants=audiobook-downloader.timer
Wants=audiobook-scheduler.service

# Translation monitor (v8.3.9) — two-tier stuck-claim / retry-budget watchdog.
# Live tier fires every 30s (origin='live', priority 0/1); sampler tier fires
# every 5min (sampler_jobs + origin in ('sampler','backlog')). See
# docs/TRANSLATION-MONITOR.md.
Wants=audiobook-translation-monitor-live.timer
Wants=audiobook-translation-monitor-sampler.timer

[Install]
WantedBy=multi-user.target
