# systemd/audiobook-periodicals-sync.timer
[Unit]
Description=Sync Audible periodicals index every 10 minutes
Documentation=file:///opt/audiobooks/docs/PERIODICALS.md

[Timer]
# Run every 10 minutes (6x per hour) to catch new subscribed content
OnCalendar=*:0/10

# Catch up if system was off during scheduled time
Persistent=true

# Add random delay up to 30 seconds to avoid API hammering
RandomizedDelaySec=30

# Don't wake system from sleep just for this
WakeSystem=false

[Install]
WantedBy=audiobook.target timers.target
