# systemd/audiobook-api.service
[Unit]
Description=Audiobooks Library API Server (Waitress)
Documentation=https://github.com/TheBoscoClub/Audiobook-Manager
After=network-online.target
Wants=network-online.target
# Ensure application and data directories are mounted before starting.
# ReadWritePaths in [Service] requires these paths to exist for namespace setup.
# If your data dir is on a separate mount (e.g., RAID array), update /srv/audiobooks
# to match — otherwise the service will fail with NAMESPACE errors until the mount is ready.
RequiresMountsFor=/opt/audiobooks /srv/audiobooks

[Service]
Type=simple
User=audiobooks
Group=audiobooks
WorkingDirectory=/opt/audiobooks/library/backend

# Load configuration from /etc/audiobooks/audiobooks.conf
EnvironmentFile=/etc/audiobooks/audiobooks.conf
Environment=PYTHONUNBUFFERED=1
# HOME is needed for audible-cli to find its config at ~/.audible/
Environment=HOME=/var/lib/audiobooks

# Only start if port is not already in use (uses ss from iproute2, always available)
ExecStartPre=/bin/sh -c '! /usr/bin/ss -tlnp | /usr/bin/grep -q ":${AUDIOBOOKS_API_PORT} "'
ExecStart=/opt/audiobooks/library/venv/bin/python api_server.py

Restart=on-failure
RestartSec=5
# Prevent rapid restart loops from flooding the journal (max 5 starts per 60s)
StartLimitIntervalSec=60
StartLimitBurst=5
# Use journal logging for better log management (query with: journalctl -u audiobook-api)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audiobook-api

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
# Writable paths - database state + data directory + /tmp for temp files
# NOTE: Systemd doesn't expand env vars in ReadWritePaths. If your data directory
# differs from the default /srv/audiobooks, add it here explicitly.
ReadWritePaths=/var/lib/audiobooks /srv/audiobooks /tmp
ProtectHome=yes
# PrivateTmp=no - API needs access to /tmp/audiobook-staging (shared with converter)
PrivateTmp=no

[Install]
WantedBy=audiobook.target multi-user.target
