# systemd/audiobook-translate-check.timer
[Unit]
Description=Periodically check for pending audiobook translations
Documentation=https://github.com/TheBoscoClub/Audiobook-Manager
PartOf=audiobook.target

[Timer]
# Check every 5 minutes for new translations to process.
# The check itself is trivial (one sqlite3 query + systemctl check).
# GPU teardown is handled INSIDE the daemon on queue empty — this timer
# only controls how quickly new content starts processing after the
# daemon has already shut down.
OnCalendar=*:0/5
# Also check 2 minutes after boot
OnBootSec=2min
Persistent=true

[Install]
WantedBy=audiobook.target
